# apps/server/systemd/vibesensor-hotspot.service
[Unit]
Description=Configure VibeSensor hotspot
After=NetworkManager.service
Wants=NetworkManager.service
After=vibesensor-rfkill-unblock.service
Wants=vibesensor-rfkill-unblock.service

[Service]
Type=oneshot
Environment=VIBESENSOR_CONFIG_PATH=/etc/vibesensor/config.yaml
ExecStart=__PI_DIR__/scripts/hotspot_nmcli.sh --config /etc/vibesensor/config.yaml
TimeoutStartSec=300
Restart=on-failure
RestartSec=5
RemainAfterExit=yes
SyslogIdentifier=vibesensor-hotspot

[Install]
WantedBy=multi-user.target
