# apps/server/systemd/vibesensor-rfkill-unblock.service
[Unit]
Description=Unblock Wi-Fi and Bluetooth radios before networking starts
After=systemd-rfkill.service
Before=bluetooth.service NetworkManager.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'if command -v rfkill >/dev/null 2>&1; then rfkill unblock bluetooth || true; rfkill unblock wifi || rfkill unblock all || true; fi'
RemainAfterExit=yes

[Install]
# Tie this preflight to both Bluetooth and NetworkManager startup so radio state
# is fixed before either service begins work.
WantedBy=bluetooth.service
WantedBy=NetworkManager.service
