# apps/server/systemd/vibesensor-rfkill-unblock.service
[Unit]
Description=Unblock Wi-Fi radios before NetworkManager starts
After=systemd-rfkill.service
Before=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'if command -v rfkill >/dev/null 2>&1; then rfkill unblock wifi || rfkill unblock all || true; fi'
RemainAfterExit=yes

[Install]
# Tie this preflight to NetworkManager startup so boot ordering stays explicit.
WantedBy=NetworkManager.service
