# pi/systemd/vibesenser.service
[Unit]
Description=VibeSenser FastAPI Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__PI_DIR__
ExecStart=__VENV_DIR__/bin/python -m vibesenser.app --config /etc/vibesenser/config.yaml
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
