# image/pi-gen/assets/vibesensor-hotspot.service
[Unit]
Description=Configure VibeSensor hotspot
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/opt/VibeSensor/pi/scripts/hotspot_nmcli.sh
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
SyslogIdentifier=vibesensor-hotspot

[Install]
WantedBy=multi-user.target
